# cf-tunnel.service
[Unit]
Description=Cloudflare Tunnel Auto-Provisioner
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-networks/
After=network-online.target
Wants=network-online.target
# Ensure SSH is running before we expose it
After=ssh.service sshd.service
# Restart when network target is restarted
PartOf=network-online.target
# Also bind to NetworkManager if available
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=exec
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/cf-tunnel-provisioner.sh
Restart=always
RestartSec=5

# Faster restart on network issues
TimeoutStartSec=120
TimeoutStopSec=30

# Security hardening
User=root
Group=root
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Allow writes to state directory and cloudflared config
ReadWritePaths=/var/lib/cf-tunnel /etc/cloudflared

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cf-tunnel

# Watchdog - restart if unresponsive
WatchdogSec=300

# Rate limiting restarts
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
